# target/linux/ramips/mt7620/base-files/etc/uci-defaults/99_set-root-password
#!/bin/sh
# 设置 root 用户默认密码为 admin
echo "root:admin" | chpasswd
exit 0
